# deploy/gunicorn.socket
[Unit]
Description=gunicorn socket for tienda_heilyn

[Socket]
ListenStream=/run/gunicorn-tienda.sock

[Install]
WantedBy=sockets.target
